# extra/katello-tracer-upload.cron
# Send a new Tracer report after a reboot
@reboot root sleep 60 && PATH=/usr/sbin:/usr/bin:/sbin:/bin /sbin/katello-tracer-upload >/dev/null 2>&1
