# extra/katello-agent-send.cron
# Send a new Tracer report after a reboot
@reboot root /sbin/katello-tracer-upload > /dev/null 2>&1
